# webservice/gtopt-webservice.service
[Unit]
Description=gtopt Web Service
Documentation=https://github.com/marcelomatus/gtopt
After=network.target

[Service]
Type=simple
User=gtopt
Group=gtopt
WorkingDirectory=/usr/local/share/gtopt/webservice

# Environment variables
Environment="PORT=3000"
Environment="GTOPT_BIN=/usr/local/bin/gtopt"
Environment="GTOPT_DATA_DIR=/var/lib/gtopt/data"
Environment="GTOPT_LOG_DIR=/var/log/gtopt"
Environment="NODE_ENV=production"

# Start the service using the installed launcher
ExecStart=/usr/local/bin/gtopt_websrv

# Restart policy
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gtopt-webservice

[Install]
WantedBy=multi-user.target
